# systemd/system/kernel-nvme-debug.service
[Unit]
Description=NVMe debug

[Service]
ExecStart=/usr/local/bin/kernel-nvme-debug nvme

[Install]
WantedBy=multi-user.target
